# service/gpio-controller.service
[Unit]
Description=Startup

[Service]
ExecStart=/home/ubuntu/wall-controller/shell/gpio-process.sh

[Install]
WantedBy=default.target