# service/wall-gateway.service
[Unit]
Description=Startup

[Service]
ExecStart=/home/ubuntu/wall-controller/service/gateway-process.sh

[Install]
WantedBy=default.target