# service/wall-controller.service
[Unit]
Description=Startup

[Service]
ExecStart=/home/ubuntu/wall-controller/shell/wall-process.sh

[Install]
WantedBy=default.target