# service/wall-gpio.service
[Unit]
Description=Startup

[Service]
ExecStart=/home/ubuntu/wall-controller/service/gpio-process.sh

[Install]
WantedBy=default.target