# service/wall-gateway.service
[Unit]
Description=GATEWAY process service

[Service]
User=ubuntu
ExecStart=/home/ubuntu/service/wall-gateway.sh

[Install]
WantedBy=default.target
